# supervisor/supervisord.conf
[supervisord]
user=root
logfile=/tmp/supervisord.log
chmod=777

[program:flask]
command="python run.py"
directory=/flask/run.py
autostart=true
autorestart=true
