# supervisor/supervisord.conf
[supervisord]
user=nobody
logfile=/supervisord.log

[program:run.py]
command="python run.py"
directory=/flask/run.py
autostart=true
autorestart=true
