# flask/supervisord.conf
[supervisord]
user=root
logfile=/tmp/supervisord.log
nodaemon=true

[program:flask]
command=/bin/bash -c "python3 run.py"
autostart=true
autorestart=true
startsecs=10
